# config/supervisord.conf
[supervisord]
logfile = /dev/null
loglevel = info
pidfile = /var/run/supervisord.pid
nodaemon = true
user = root

[program:service-one]
command=python -u app.py
directory=/usr/src/service-one/src
stdout_logfile=/usr/src/log/service-one.txt
redirect_stderr=true
autorestart=true

[program:service-two]
command=python -u app.py
directory=/usr/src/service-two/src
stdout_logfile=/usr/src/log/service-two.txt
redirect_stderr=true
autorestart=true

[inet_http_server]
port=*:9001
username=testuser
password=testpass